# systemd/robot-brain.service.example
[Unit]
Description=Reachy Mini Agent — main brain (vision + voice + LLM)
After=network-online.target kokoro.service whisper.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/reachy-mini-agent
ExecStart=%h/reachy-mini-agent/scripts/run_robot.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/reachy-mini-agent/logs/robot_brain.log
StandardError=inherit

[Install]
WantedBy=default.target
